# services/nscd/run
#!/bin/bash
export LANG=C
export LC_ALL=

mkdir -p /var/run/nscd/

exec /usr/sbin/nscd -d >/dev/null
